# roles/docker-on-cumulus/files/docker-system-prune.service
[Unit]
Description=Docker System Prune

[Service]
Type=simple
ExecStart=/usr/bin/docker system prune -a -f